# .vim/templates/skeleton.timer
[Unit]
Description=<++>

[Timer]
OnCalendar=DayOfWeek Year-Month-Day Hour:Minute:Second
Persistant=true

[Install]
WantedBy=timers.target
