# .vim/templates/skeleton.service
[Unit]
Description=<++>
After=network.target

[Service]
Type=simple
ExecStart=<++>
ExecStop=killall <++>

[Install]
WantedBy=multi-user.target
